# deployment/api-gateway.service
[Unit]
Description=Samyama API Gateway
After=network.target samyama.service
Requires=samyama.service

[Service]
Type=simple
User=azureuser
WorkingDirectory=/home/azureuser/api-gateway
ExecStart=/usr/bin/node /home/azureuser/api-gateway/index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=8080

[Install]
WantedBy=multi-user.target
